# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:serfagent]
command=/usr/local/bin/serf.sh
